# rootfs/etc/services.d/01-glusterfs/run
#!/usr/bin/execlineb -P
with-contenv
/usr/local/sbin/glusterd 
    --log-level INFO 
    --log-file /proc/1/fd/1
    -p /var/run/glusterd.pid
    -N
