# rootfs/etc/services.d/02-rpc.idmapd/run
#!/usr/bin/execlineb -P
rpc.idmapd -f -v
