# rootfs/etc/services.d/02-rpcbind/run
#!/usr/bin/execlineb -P
s6-fghack rpcbind
