# rootfs/etc/services.d/glusterfs/run
#!/usr/bin/execlineb -P
with-contenv
/usr/local/sbin/glusterd -N
