# rootfs/etc/services.d/ganesha/run
#!/usr/bin/execlineb -P
ganesha.nfsd -F -L /dev/stdout
