# rootfs/etc/services.d/03-ganesha/run
#!/usr/bin/execlineb -P
with-contenv
/usr/bin/ganesha.nfsd -F -L /dev/stdout -N NIV_EVENT
