# rootfs/etc/services.d/01-glustereventsd/run
#!/usr/bin/execlineb -P
with-contenv
/usr/local/sbin/glustereventsd
    -p /var/run/glustereventsd.pid
